# Sauron_the_all_seeing_eye/files/crontab
@reboot root sh /home/varun/code/z_archive/autorun/Sauron_the_all_seeing_eye/sauron.sh
